# test_supervisor/myconfigfile.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[supervisord]
logfile=/tmp/supervisord.log ; main log file; default $CWD/supervisord.log
logfile_maxbytes=50MB        ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10           ; # of main logfile backups; 0 means none, default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
pidfile=/tmp/supervisord.pid ; supervisord pidfile; default supervisord.pid
nodaemon=false               ; start in foreground if true; default false
silent=false                 ; no logs to stdout if true; default false
minfds=1024                  ; min. avail startup file descriptors; default 1024
minprocs=200                 ; min. avail process descriptors;default 200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:sleep_loop]
command=/bin/sleep 3600
autostart=true
numprocs=10
process_name=sleep_loop_%(process_num)s
autorestart=true
stderr_logfile=/tmp/sleep_loop.err.log
stdout_logfile=/tmp/sleep_loop.out.log
stopsignal=HUP
stoptime=20

[program:ping_google]
command=/bin/sleep 3
autostart=true
startsecs=1
startretries=3
autorestart=true
stderr_logfile=/tmp/ping_google.err.log

[program:taiv_null]
command=/usr/bin/tail -f /dev/null
autostart=true
autorestart=true
stderr_logfile=/tmp/tail_dev_null.err.log
stdout_logfile=/tmp/t

[program:unstoppable]
command=/bin/bash -c "trap '' TERM; echo 'Signal TERM ignoré, en attente...'; while true; do /bin/sleep 1; done"

; --- Contrôle du processus ---
autostart=true
stopsignal=TERM
stopwaitsecs=5  ; Après 5s d'attente, supervisor enverra SIGKILL

; --- Fichiers de log ---
stdout_logfile=/tmp/unstoppable.out
stderr_logfile=/tmp/unstoppable.err

[program:echo_loop]
command=/bin/bash -c "while true; do echo 'Hello from Supervisor'; sleep 10; done"
autostart=true
autorestart=true
stderr_logfile=/tmp/echo_loop.err.log
stdout_logfile=/tmp/echo_loop.out.log

[program:test_stop_signal]
command=/bin/bash /app/test_supervisor/programmes/ignore_SIGHUP.sh
autostart=true
autorestart=true
stderr_logfile=/tmp/test_stop_signal.err.log
stdout_logfile=/tmp/test_stop_signal.out.log
stopsignal=HUP
stoptime=20

[program:test_stop_signal2]
command=/bin/bash /app/test_supervisor/programmes/ignore_SIGHUP.sh
autostart=true
autorestart=true
stderr_logfile=/tmp/test_stop_signal2.err.log
stdout_logfile=/tmp/test_stop_signal2.out.log
stopsignal=HUP
stoptime=20